# ops/cron/noctification-db-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
APP_ROOT=/home/noctification/noctification
DB_FILE=/home/noctification/noctification/apps/api/data/noctification.db
BACKUP_DIR=/home/noctification/noctification/backups/db
RETENTION_DAYS=14

0 2 * * * noctification $APP_ROOT/ops/scripts/backup-db.sh >> $APP_ROOT/backups/db/noctification-backup.log 2>&1
