# ops/cron/noctification-db-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
APP_ROOT=__APP_ROOT__
DB_FILE=__APP_ROOT__/apps/api/data/noctification.db
BACKUP_DIR=__APP_ROOT__/backups/db
RETENTION_DAYS=14

0 2 * * * __SYSTEM_USER__ $APP_ROOT/ops/scripts/backup-db.sh >> $APP_ROOT/backups/db/noctification-backup.log 2>&1
