# ops/cron/noctification-db-backup.cron
0 2 * * * DB_FILE=/opt/noctification/apps/api/data/noctification.db BACKUP_DIR=/opt/noctification/backups/db RETENTION_DAYS=14 /opt/noctification/ops/scripts/backup-db.sh >> /var/log/noctification-backup.log 2>&1
