# ops/systemd/noctification-api.user.service
[Unit]
Description=Noctification API (user)
After=network.target

[Service]
Type=simple
Environment=APP_ROOT=%h/noctification
EnvironmentFile=-%h/.config/noctification/api.env
ExecStart=/usr/bin/bash -lc 'cd "${APP_ROOT}/apps/api" && exec /usr/bin/node dist/index.js'
Restart=always
RestartSec=5
KillSignal=SIGINT

[Install]
WantedBy=default.target
