# ops/systemd/noctification-api.service
[Unit]
Description=Noctification API
After=network.target

[Service]
Type=simple
User=noctification
Group=noctification
Environment=APP_ROOT=/opt/noctification
EnvironmentFile=/etc/noctification/api.env
ExecStart=/usr/bin/bash -lc 'cd "${APP_ROOT}/apps/api" && exec /usr/bin/node dist/index.js'
Restart=always
RestartSec=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
